# packages/reorg-monitor/reorg-monitor/etc/cron.daily/dpkg
#!/bin/sh

# Skip if systemd is running.
if [ -d /run/systemd/system ]; then
  exit 0
fi

/usr/libexec/dpkg/dpkg-db-backup
