# packages/reorg-monitor/reorg-monitor/var/lib/dpkg/info/apt.conffiles
/etc/apt/apt.conf.d/01autoremove
/etc/cron.daily/apt-compat
/etc/logrotate.d/apt
