# packages/reorg-monitor/reorg-monitor/var/lib/dpkg/info/apt.list
/.
/etc
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/01autoremove
/etc/apt/auth.conf.d
/etc/apt/keyrings
/etc/apt/preferences.d
/etc/apt/sources.list.d
/etc/apt/trusted.gpg.d
/etc/cron.daily
/etc/cron.daily/apt-compat
/etc/kernel
/etc/kernel/postinst.d
/etc/logrotate.d
/etc/logrotate.d/apt
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/apt-daily-upgrade.service
/lib/systemd/system/apt-daily-upgrade.timer
/lib/systemd/system/apt-daily.service
/lib/systemd/system/apt-daily.timer
/usr
/usr/bin
/usr/bin/apt
/usr/bin/apt-cache
/usr/bin/apt-cdrom
/usr/bin/apt-config
/usr/bin/apt-get
/usr/bin/apt-key
/usr/bin/apt-mark
/usr/lib
/usr/lib/aarch64-linux-gnu
/usr/lib/aarch64-linux-gnu/libapt-private.so.0.0.0
/usr/lib/apt
/usr/lib/apt/apt-helper
/usr/lib/apt/apt.systemd.daily
/usr/lib/apt/methods
/usr/lib/apt/methods/cdrom
/usr/lib/apt/methods/copy
/usr/lib/apt/methods/file
/usr/lib/apt/methods/ftp
/usr/lib/apt/methods/gpgv
/usr/lib/apt/methods/http
/usr/lib/apt/methods/mirror
/usr/lib/apt/methods/rred
/usr/lib/apt/methods/rsh
/usr/lib/apt/methods/store
/usr/lib/apt/planners
/usr/lib/apt/solvers
/usr/lib/apt/solvers/dump
/usr/lib/dpkg
/usr/lib/dpkg/methods
/usr/lib/dpkg/methods/apt
/usr/lib/dpkg/methods/apt/desc.apt
/usr/lib/dpkg/methods/apt/install
/usr/lib/dpkg/methods/apt/names
/usr/lib/dpkg/methods/apt/setup
/usr/lib/dpkg/methods/apt/update
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/apt
/usr/share/bug
/usr/share/bug/apt
/usr/share/bug/apt/script
/usr/share/doc
/usr/share/doc/apt
/usr/share/doc/apt/NEWS.Debian.gz
/usr/share/doc/apt/README.md.gz
/usr/share/doc/apt/changelog.gz
/usr/share/doc/apt/copyright
/usr/share/doc/apt/examples
/usr/share/doc/apt/examples/apt.conf
/usr/share/doc/apt/examples/configure-index
/usr/share/doc/apt/examples/preferences
/usr/share/doc/apt/examples/sources.list
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/apt
/usr/share/locale
/usr/share/locale/ar
/usr/share/locale/ar/LC_MESSAGES
/usr/share/locale/ar/LC_MESSAGES/apt.mo
/usr/share/locale/ast
/usr/share/locale/ast/LC_MESSAGES
/usr/share/locale/ast/LC_MESSAGES/apt.mo
/usr/share/locale/bg
/usr/share/locale/bg/LC_MESSAGES
/usr/share/locale/bg/LC_MESSAGES/apt.mo
/usr/share/locale/bs
/usr/share/locale/bs/LC_MESSAGES
/usr/share/locale/bs/LC_MESSAGES/apt.mo
/usr/share/locale/ca
/usr/share/locale/ca/LC_MESSAGES
/usr/share/locale/ca/LC_MESSAGES/apt.mo
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/cs/LC_MESSAGES/apt.mo
/usr/share/locale/cy
/usr/share/locale/cy/LC_MESSAGES
/usr/share/locale/cy/LC_MESSAGES/apt.mo
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/apt.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/apt.mo
/usr/share/locale/dz
/usr/share/locale/dz/LC_MESSAGES
/usr/share/locale/dz/LC_MESSAGES/apt.mo
/usr/share/locale/el
/usr/share/locale/el/LC_MESSAGES
/usr/share/locale/el/LC_MESSAGES/apt.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/apt.mo
/usr/share/locale/eu
/usr/share/locale/eu/LC_MESSAGES
/usr/share/locale/eu/LC_MESSAGES/apt.mo
/usr/share/locale/fi
/usr/share/locale/fi/LC_MESSAGES
/usr/share/locale/fi/LC_MESSAGES/apt.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/apt.mo
/usr/share/locale/gl
/usr/share/locale/gl/LC_MESSAGES
/usr/share/locale/gl/LC_MESSAGES/apt.mo
/usr/share/locale/hu
/usr/share/locale/hu/LC_MESSAGES
/usr/share/locale/hu/LC_MESSAGES/apt.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/apt.mo
/usr/share/locale/ja
/usr/share/locale/ja/LC_MESSAGES
/usr/share/locale/ja/LC_MESSAGES/apt.mo
/usr/share/locale/km
/usr/share/locale/km/LC_MESSAGES
/usr/share/locale/km/LC_MESSAGES/apt.mo
/usr/share/locale/ko
/usr/share/locale/ko/LC_MESSAGES
/usr/share/locale/ko/LC_MESSAGES/apt.mo
/usr/share/locale/ku
/usr/share/locale/ku/LC_MESSAGES
/usr/share/locale/ku/LC_MESSAGES/apt.mo
/usr/share/locale/lt
/usr/share/locale/lt/LC_MESSAGES
/usr/share/locale/lt/LC_MESSAGES/apt.mo
/usr/share/locale/mr
/usr/share/locale/mr/LC_MESSAGES
/usr/share/locale/mr/LC_MESSAGES/apt.mo
/usr/share/locale/nb
/usr/share/locale/nb/LC_MESSAGES
/usr/share/locale/nb/LC_MESSAGES/apt.mo
/usr/share/locale/ne
/usr/share/locale/ne/LC_MESSAGES
/usr/share/locale/ne/LC_MESSAGES/apt.mo
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/nl/LC_MESSAGES/apt.mo
/usr/share/locale/nn
/usr/share/locale/nn/LC_MESSAGES
/usr/share/locale/nn/LC_MESSAGES/apt.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/apt.mo
/usr/share/locale/pt
/usr/share/locale/pt/LC_MESSAGES
/usr/share/locale/pt/LC_MESSAGES/apt.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/apt.mo
/usr/share/locale/ro
/usr/share/locale/ro/LC_MESSAGES
/usr/share/locale/ro/LC_MESSAGES/apt.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/apt.mo
/usr/share/locale/sk
/usr/share/locale/sk/LC_MESSAGES
/usr/share/locale/sk/LC_MESSAGES/apt.mo
/usr/share/locale/sl
/usr/share/locale/sl/LC_MESSAGES
/usr/share/locale/sl/LC_MESSAGES/apt.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/apt.mo
/usr/share/locale/th
/usr/share/locale/th/LC_MESSAGES
/usr/share/locale/th/LC_MESSAGES/apt.mo
/usr/share/locale/tl
/usr/share/locale/tl/LC_MESSAGES
/usr/share/locale/tl/LC_MESSAGES/apt.mo
/usr/share/locale/tr
/usr/share/locale/tr/LC_MESSAGES
/usr/share/locale/tr/LC_MESSAGES/apt.mo
/usr/share/locale/uk
/usr/share/locale/uk/LC_MESSAGES
/usr/share/locale/uk/LC_MESSAGES/apt.mo
/usr/share/locale/vi
/usr/share/locale/vi/LC_MESSAGES
/usr/share/locale/vi/LC_MESSAGES/apt.mo
/usr/share/locale/zh_CN
/usr/share/locale/zh_CN/LC_MESSAGES
/usr/share/locale/zh_CN/LC_MESSAGES/apt.mo
/usr/share/locale/zh_TW
/usr/share/locale/zh_TW/LC_MESSAGES
/usr/share/locale/zh_TW/LC_MESSAGES/apt.mo
/usr/share/man
/usr/share/man/de
/usr/share/man/de/man1
/usr/share/man/de/man1/apt-transport-http.1.gz
/usr/share/man/de/man1/apt-transport-https.1.gz
/usr/share/man/de/man1/apt-transport-mirror.1.gz
/usr/share/man/de/man5
/usr/share/man/de/man5/apt.conf.5.gz
/usr/share/man/de/man5/apt_auth.conf.5.gz
/usr/share/man/de/man5/apt_preferences.5.gz
/usr/share/man/de/man5/sources.list.5.gz
/usr/share/man/de/man7
/usr/share/man/de/man7/apt-patterns.7.gz
/usr/share/man/de/man8
/usr/share/man/de/man8/apt-cache.8.gz
/usr/share/man/de/man8/apt-cdrom.8.gz
/usr/share/man/de/man8/apt-config.8.gz
/usr/share/man/de/man8/apt-get.8.gz
/usr/share/man/de/man8/apt-mark.8.gz
/usr/share/man/de/man8/apt-secure.8.gz
/usr/share/man/de/man8/apt.8.gz
/usr/share/man/es
/usr/share/man/es/man5
/usr/share/man/es/man5/apt_preferences.5.gz
/usr/share/man/es/man8
/usr/share/man/es/man8/apt-cache.8.gz
/usr/share/man/es/man8/apt-cdrom.8.gz
/usr/share/man/es/man8/apt-config.8.gz
/usr/share/man/fr
/usr/share/man/fr/man1
/usr/share/man/fr/man1/apt-transport-http.1.gz
/usr/share/man/fr/man1/apt-transport-https.1.gz
/usr/share/man/fr/man1/apt-transport-mirror.1.gz
/usr/share/man/fr/man5
/usr/share/man/fr/man5/apt.conf.5.gz
/usr/share/man/fr/man5/apt_auth.conf.5.gz
/usr/share/man/fr/man5/apt_preferences.5.gz
/usr/share/man/fr/man5/sources.list.5.gz
/usr/share/man/fr/man8
/usr/share/man/fr/man8/apt-cache.8.gz
/usr/share/man/fr/man8/apt-cdrom.8.gz
/usr/share/man/fr/man8/apt-config.8.gz
/usr/share/man/fr/man8/apt-get.8.gz
/usr/share/man/fr/man8/apt-mark.8.gz
/usr/share/man/fr/man8/apt-secure.8.gz
/usr/share/man/fr/man8/apt.8.gz
/usr/share/man/it
/usr/share/man/it/man5
/usr/share/man/it/man5/apt.conf.5.gz
/usr/share/man/it/man5/apt_preferences.5.gz
/usr/share/man/it/man5/sources.list.5.gz
/usr/share/man/it/man8
/usr/share/man/it/man8/apt-cache.8.gz
/usr/share/man/it/man8/apt-cdrom.8.gz
/usr/share/man/it/man8/apt-config.8.gz
/usr/share/man/it/man8/apt-get.8.gz
/usr/share/man/it/man8/apt-mark.8.gz
/usr/share/man/it/man8/apt-secure.8.gz
/usr/share/man/it/man8/apt.8.gz
/usr/share/man/ja
/usr/share/man/ja/man5
/usr/share/man/ja/man5/apt.conf.5.gz
/usr/share/man/ja/man5/apt_preferences.5.gz
/usr/share/man/ja/man5/sources.list.5.gz
/usr/share/man/ja/man8
/usr/share/man/ja/man8/apt-cache.8.gz
/usr/share/man/ja/man8/apt-cdrom.8.gz
/usr/share/man/ja/man8/apt-config.8.gz
/usr/share/man/ja/man8/apt-get.8.gz
/usr/share/man/ja/man8/apt-mark.8.gz
/usr/share/man/ja/man8/apt-secure.8.gz
/usr/share/man/ja/man8/apt.8.gz
/usr/share/man/man1
/usr/share/man/man1/apt-transport-http.1.gz
/usr/share/man/man1/apt-transport-https.1.gz
/usr/share/man/man1/apt-transport-mirror.1.gz
/usr/share/man/man5
/usr/share/man/man5/apt.conf.5.gz
/usr/share/man/man5/apt_auth.conf.5.gz
/usr/share/man/man5/apt_preferences.5.gz
/usr/share/man/man5/sources.list.5.gz
/usr/share/man/man7
/usr/share/man/man7/apt-patterns.7.gz
/usr/share/man/man8
/usr/share/man/man8/apt-cache.8.gz
/usr/share/man/man8/apt-cdrom.8.gz
/usr/share/man/man8/apt-config.8.gz
/usr/share/man/man8/apt-get.8.gz
/usr/share/man/man8/apt-key.8.gz
/usr/share/man/man8/apt-mark.8.gz
/usr/share/man/man8/apt-secure.8.gz
/usr/share/man/man8/apt.8.gz
/usr/share/man/nl
/usr/share/man/nl/man1
/usr/share/man/nl/man1/apt-transport-http.1.gz
/usr/share/man/nl/man1/apt-transport-https.1.gz
/usr/share/man/nl/man1/apt-transport-mirror.1.gz
/usr/share/man/nl/man5
/usr/share/man/nl/man5/apt.conf.5.gz
/usr/share/man/nl/man5/apt_auth.conf.5.gz
/usr/share/man/nl/man5/apt_preferences.5.gz
/usr/share/man/nl/man5/sources.list.5.gz
/usr/share/man/nl/man7
/usr/share/man/nl/man7/apt-patterns.7.gz
/usr/share/man/nl/man8
/usr/share/man/nl/man8/apt-cache.8.gz
/usr/share/man/nl/man8/apt-cdrom.8.gz
/usr/share/man/nl/man8/apt-config.8.gz
/usr/share/man/nl/man8/apt-get.8.gz
/usr/share/man/nl/man8/apt-key.8.gz
/usr/share/man/nl/man8/apt-mark.8.gz
/usr/share/man/nl/man8/apt-secure.8.gz
/usr/share/man/nl/man8/apt.8.gz
/usr/share/man/pl
/usr/share/man/pl/man5
/usr/share/man/pl/man5/apt_preferences.5.gz
/usr/share/man/pl/man8
/usr/share/man/pl/man8/apt-cache.8.gz
/usr/share/man/pl/man8/apt-cdrom.8.gz
/usr/share/man/pl/man8/apt-config.8.gz
/usr/share/man/pt
/usr/share/man/pt/man1
/usr/share/man/pt/man1/apt-transport-http.1.gz
/usr/share/man/pt/man1/apt-transport-https.1.gz
/usr/share/man/pt/man1/apt-transport-mirror.1.gz
/usr/share/man/pt/man5
/usr/share/man/pt/man5/apt.conf.5.gz
/usr/share/man/pt/man5/apt_auth.conf.5.gz
/usr/share/man/pt/man5/apt_preferences.5.gz
/usr/share/man/pt/man5/sources.list.5.gz
/usr/share/man/pt/man7
/usr/share/man/pt/man7/apt-patterns.7.gz
/usr/share/man/pt/man8
/usr/share/man/pt/man8/apt-cache.8.gz
/usr/share/man/pt/man8/apt-cdrom.8.gz
/usr/share/man/pt/man8/apt-config.8.gz
/usr/share/man/pt/man8/apt-get.8.gz
/usr/share/man/pt/man8/apt-key.8.gz
/usr/share/man/pt/man8/apt-mark.8.gz
/usr/share/man/pt/man8/apt-secure.8.gz
/usr/share/man/pt/man8/apt.8.gz
/var
/var/cache
/var/cache/apt
/var/cache/apt/archives
/var/cache/apt/archives/partial
/var/lib
/var/lib/apt
/var/lib/apt/lists
/var/lib/apt/lists/partial
/var/lib/apt/mirrors
/var/lib/apt/mirrors/partial
/var/lib/apt/periodic
/var/log
/var/log/apt
/usr/lib/aarch64-linux-gnu/libapt-private.so.0.0
/usr/lib/apt/methods/https
/usr/lib/apt/methods/mirror+copy
/usr/lib/apt/methods/mirror+file
/usr/lib/apt/methods/mirror+ftp
/usr/lib/apt/methods/mirror+http
/usr/lib/apt/methods/mirror+https
/usr/lib/apt/methods/ssh
/usr/lib/apt/planners/dump
