# packages/reorg-monitor/reorg-monitor/var/lib/dpkg/info/dpkg.conffiles
/etc/alternatives/README
/etc/cron.daily/dpkg
/etc/dpkg/dpkg.cfg
/etc/logrotate.d/alternatives
/etc/logrotate.d/dpkg
