# packages/reorg-monitor/reorg-monitor/var/lib/dpkg/info/e2fsprogs.list
/.
/etc
/etc/cron.d
/etc/cron.d/e2scrub_all
/etc/e2scrub.conf
/etc/mke2fs.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/e2scrub@.service
/lib/systemd/system/e2scrub_all.service
/lib/systemd/system/e2scrub_all.timer
/lib/systemd/system/e2scrub_fail@.service
/lib/systemd/system/e2scrub_reap.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/96-e2scrub.rules
/sbin
/sbin/badblocks
/sbin/debugfs
/sbin/dumpe2fs
/sbin/e2fsck
/sbin/e2image
/sbin/e2scrub
/sbin/e2scrub_all
/sbin/e2undo
/sbin/mke2fs
/sbin/resize2fs
/sbin/tune2fs
/usr
/usr/bin
/usr/bin/chattr
/usr/bin/lsattr
/usr/lib
/usr/lib/aarch64-linux-gnu
/usr/lib/aarch64-linux-gnu/e2fsprogs
/usr/lib/aarch64-linux-gnu/e2fsprogs/e2scrub_all_cron
/usr/lib/aarch64-linux-gnu/e2fsprogs/e2scrub_fail
/usr/sbin
/usr/sbin/e2freefrag
/usr/sbin/e4crypt
/usr/sbin/e4defrag
/usr/sbin/filefrag
/usr/sbin/mklost+found
/usr/share
/usr/share/doc
/usr/share/doc/e2fsprogs
/usr/share/doc/e2fsprogs/NEWS.gz
/usr/share/doc/e2fsprogs/README
/usr/share/doc/e2fsprogs/changelog.Debian.arm64.gz
/usr/share/doc/e2fsprogs/changelog.Debian.gz
/usr/share/doc/e2fsprogs/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/e2fsprogs
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/chattr.1.gz
/usr/share/man/man1/lsattr.1.gz
/usr/share/man/man5
/usr/share/man/man5/e2fsck.conf.5.gz
/usr/share/man/man5/ext4.5.gz
/usr/share/man/man5/mke2fs.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/badblocks.8.gz
/usr/share/man/man8/debugfs.8.gz
/usr/share/man/man8/dumpe2fs.8.gz
/usr/share/man/man8/e2freefrag.8.gz
/usr/share/man/man8/e2fsck.8.gz
/usr/share/man/man8/e2image.8.gz
/usr/share/man/man8/e2label.8.gz
/usr/share/man/man8/e2mmpstatus.8.gz
/usr/share/man/man8/e2scrub.8.gz
/usr/share/man/man8/e2scrub_all.8.gz
/usr/share/man/man8/e2undo.8.gz
/usr/share/man/man8/e4crypt.8.gz
/usr/share/man/man8/e4defrag.8.gz
/usr/share/man/man8/filefrag.8.gz
/usr/share/man/man8/mke2fs.8.gz
/usr/share/man/man8/mklost+found.8.gz
/usr/share/man/man8/resize2fs.8.gz
/usr/share/man/man8/tune2fs.8.gz
/sbin/e2label
/sbin/e2mmpstatus
/sbin/fsck.ext2
/sbin/fsck.ext3
/sbin/fsck.ext4
/sbin/mkfs.ext2
/sbin/mkfs.ext3
/sbin/mkfs.ext4
/usr/share/man/man5/ext2.5.gz
/usr/share/man/man5/ext3.5.gz
/usr/share/man/man8/fsck.ext2.8.gz
/usr/share/man/man8/fsck.ext3.8.gz
/usr/share/man/man8/fsck.ext4.8.gz
/usr/share/man/man8/mkfs.ext2.8.gz
/usr/share/man/man8/mkfs.ext3.8.gz
/usr/share/man/man8/mkfs.ext4.8.gz
