# packages/reorg-monitor/reorg-monitor/var/lib/dpkg/info/libpam-modules-bin.list
/.
/sbin
/sbin/mkhomedir_helper
/sbin/pam_namespace_helper
/sbin/pwhistory_helper
/sbin/unix_chkpwd
/sbin/unix_update
/usr
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/pam_namespace.service
/usr/sbin
/usr/sbin/faillock
/usr/sbin/pam_timestamp_check
/usr/share
/usr/share/doc
/usr/share/doc/libpam-modules-bin
/usr/share/doc/libpam-modules-bin/NEWS.Debian.gz
/usr/share/doc/libpam-modules-bin/changelog.Debian.gz
/usr/share/doc/libpam-modules-bin/changelog.gz
/usr/share/doc/libpam-modules-bin/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/libpam-modules-bin
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/environment.5.gz
/usr/share/man/man8
/usr/share/man/man8/faillock.8.gz
/usr/share/man/man8/mkhomedir_helper.8.gz
/usr/share/man/man8/pam_namespace_helper.8.gz
/usr/share/man/man8/pam_timestamp_check.8.gz
/usr/share/man/man8/pwhistory_helper.8.gz
/usr/share/man/man8/unix_chkpwd.8.gz
/usr/share/man/man8/unix_update.8.gz
