# packages/reorg-monitor/reorg-monitor/var/lib/dpkg/info/util-linux-extra.list
/.
/etc
/etc/default
/etc/default/hwclock
/etc/init.d
/etc/init.d/hwclock.sh
/sbin
/sbin/hwclock
/usr
/usr/bin
/usr/bin/fincore
/usr/bin/lsfd
/usr/bin/lsirq
/usr/lib
/usr/lib/udev
/usr/lib/udev/hwclock-set
/usr/lib/udev/rules.d
/usr/lib/udev/rules.d/85-hwclock.rules
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/fincore
/usr/share/bash-completion/completions/hwclock
/usr/share/bash-completion/completions/lsirq
/usr/share/doc
/usr/share/doc/util-linux-extra
/usr/share/doc/util-linux-extra/changelog.Debian.gz
/usr/share/doc/util-linux-extra/changelog.gz
/usr/share/doc/util-linux-extra/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/fincore.1.gz
/usr/share/man/man1/lsfd.1.gz
/usr/share/man/man1/lsirq.1.gz
/usr/share/man/man5
/usr/share/man/man5/hwclock.5.gz
/usr/share/man/man8
/usr/share/man/man8/hwclock.8.gz
