# packages/reorg-monitor/reorg-monitor/usr/lib/systemd/system/dpkg-db-backup.service
[Unit]
Description=Daily dpkg database backup service
Documentation=man:dpkg(1)

[Service]
Type=oneshot
ExecStart=/usr/libexec/dpkg/dpkg-db-backup
