# packages/reorg-monitor/reorg-monitor/usr/lib/systemd/system/dpkg-db-backup.timer
[Unit]
Description=Daily dpkg database backup timer
Documentation=man:dpkg(1)

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
