# packages/reorg-monitor/reorg-monitor/usr/lib/systemd/system/pam_namespace.service
[Unit]
After=local-fs.target
Before=multi-user.target shutdown.target
Conflicts=shutdown.target
DefaultDependencies=no
Description=Make sure parent directories configured in /etc/security/namespace.conf for polyinstantiation exist
Documentation=man:pam_namespace(8)

[Service]
ExecStart=/sbin/pam_namespace_helper
Type=oneshot
