# packages/reorg-monitor/reorg-monitor/usr/lib/udev/rules.d/85-hwclock.rules
# Set the System Time from the Hardware Clock and set the kernel's timezone
# value to the local timezone when the kernel clock module is loaded.

KERNEL=="rtc0", RUN+="/usr/lib/udev/hwclock-set $root/$name"
